# systemd/phantom.service
[Unit]
Description=Phantom Browser Automation Agent
Documentation=https://github.com/NinjaTech-AI/phantom
After=network.target supervisor.service
Wants=supervisor.service

[Service]
Type=simple
User=root
WorkingDirectory=/workspace/phantom
ExecStartPre=/bin/sleep 10
ExecStart=/usr/local/bin/python orchestrator.py
Restart=on-failure
RestartSec=30s

Environment=HOME=/root
Environment=DISPLAY=:99
Environment=PYTHONPATH=/workspace
StandardOutput=append:/workspace/logs/phantom-service.log
StandardError=append:/workspace/logs/phantom-service.log
SyslogIdentifier=phantom

[Install]
WantedBy=multi-user.target
